# scripts/autosave.conf
description     "DevTools Autosave Server Upstart job"
author          "Nikita Vasilyev"

start on runlevel [2345]
stop on runlevel [06]

respawn

exec start-stop-daemon --start --make-pidfile --pidfile /var/run/autosave.pid --exec /usr/local/bin/autosave >> /var/log/autosave.log 2>&1
